# files/etc/crontabs/daemon
0 */3 * * * /etc/init.d/collectd restart
