# files/etc/crontabs/daemon
*/4 3 * * * /etc/init.d/collectd restart
