# files/etc/rc.local
#!/bin/sh

LEGACY=$(grep legacy /etc/freifunk_version)
if [ $LEGACY ]; then
        echo "1" > /proc/sys/net/ipv6/conf/br-wan/accept_ra
fi

