# files/etc/rc.local
#!/bin/sh
. /lib/freifunk/lib_node.sh

if is_legacy ; then
        echo "1" > /proc/sys/net/ipv6/conf/br-wan/disable_ipv6
fi

